# systemd/mhde-fx-retrain.timer
[Unit]
Description=MHDE FX weekly retrain timer

[Timer]
OnCalendar=Sat *-*-* 22:00:00
Persistent=true

[Install]
WantedBy=timers.target
