# systemd/mhde-monitor-pipeline.service
[Unit]
Description=MHDE monitor: pipeline execution check (recency + row count)
After=network.target

[Service]
Type=oneshot
User=jpcg
WorkingDirectory=/home/jpcg/MHDE
Environment=MHDE_DB_PATH=/home/jpcg/MHDE/data/mhde.duckdb
ExecStart=/home/jpcg/MHDE/venv/bin/python main.py monitor pipeline-execution
StandardOutput=append:/home/jpcg/MHDE/data/logs/monitor_pipeline.log
StandardError=append:/home/jpcg/MHDE/data/logs/monitor_pipeline.log
