# systemd/mhde-crypto-pipeline-monitor.timer
[Unit]
Description=MHDE crypto pipeline monitor — daily 00:50 UTC (after the 00:40 export and the ~00:45 engine entry)

[Timer]
OnCalendar=*-*-* 00:50:00
Persistent=true

[Install]
WantedBy=timers.target
