# systemd/mhde-monitor-config-drift.timer
[Unit]
Description=MHDE monitor: config drift daily

[Timer]
OnCalendar=*-*-* 12:15:00
Persistent=true

[Install]
WantedBy=timers.target
