# systemd/mhde-capture-core.service
# /home/jpcg/.config/systemd/user/mhde-capture-core.service
[Unit]
Description=MHDE capture-core — raw Binance USDT-M market-data capture (read-only, public)
After=network-online.target
Wants=network-online.target

[Service]
# systemd --user unit (installed under ~/.config/systemd/user/): runs as the
# logging-in user, so NO User=/Group= line. Persistence relies on linger
# (`loginctl enable-linger jpcg`). Pre-flight fails fast if the venv is missing.
# The service captures ONLY raw parquet under data/research/capture_core/ and
# never opens the production database, so no production-DB pre-flight is needed.
#
# BUILT-NOT-DEPLOYED: this file is tracked but not installed/enabled. The disk
# guard (free-space prune@soft / halt@critical, firehose-only) is enforced inside
# the capture process (ADR-036). R2 offload is DROPPED (ADR-035: forward-only).
# Resource caps below express PRIORITY, not starvation: under contention the engine
# (system.slice) wins CPU+IO, but capture is never throttled to zero.
#
# Caps are tightened at deploy per the OPERATIONS.md runbook (MemoryMax -> ~peak
# RSS x1.8 after measurement). The cross-slice arbitration that lets the engine
# win requires the SUDO-installed drop-in systemd/system/user.slice.d/
# 10-capture-deprioritize.conf; per-unit weights alone only arbitrate WITHIN
# user.slice.
Type=simple
WorkingDirectory=/home/jpcg/MHDE
ExecStartPre=/usr/bin/test -x /home/jpcg/MHDE/venv/bin/python
ExecStart=/home/jpcg/MHDE/venv/bin/python main.py crypto capture-core-run
Restart=on-failure
RestartSec=5
# Give the writer time to flush open parquet buffers on SIGTERM before SIGKILL.
TimeoutStopSec=30
# Resource caps (PR-3): MemoryMax is the firehose ceiling (dominant writer);
# CPUWeight/IOWeight < the system.slice default (100) so the engine wins under
# contention; OOMScoreAdjust makes capture the first OOM victim, never the engine.
MemoryMax=4G
CPUWeight=20
IOWeight=20
OOMScoreAdjust=800
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
