# systemd/mhde-monitor-phase0-calibration.service
[Unit]
Description=MHDE monitor: Phase 0 calibration weekly check
After=network.target

[Service]
Type=oneshot
User=jpcg
WorkingDirectory=/home/jpcg/MHDE
Environment=MHDE_DB_PATH=/home/jpcg/MHDE/data/mhde.duckdb
ExecStart=/home/jpcg/MHDE/venv/bin/python main.py monitor phase0-calibration
StandardOutput=append:/home/jpcg/MHDE/data/logs/monitor_phase0_calibration.log
StandardError=append:/home/jpcg/MHDE/data/logs/monitor_phase0_calibration.log
