# systemd/mhde-crypto-retrain.service
[Unit]
Description=MHDE crypto weekly retrain
After=network.target

[Service]
Type=oneshot
User=jpcg
WorkingDirectory=/home/jpcg/MHDE
ExecStart=/home/jpcg/MHDE/venv/bin/python main.py crypto retrain
Environment=MHDE_DB_PATH=/home/jpcg/MHDE/data/mhde.duckdb
StandardOutput=append:/home/jpcg/MHDE/data/logs/crypto_retrain.log
StandardError=append:/home/jpcg/MHDE/data/logs/crypto_retrain.log
TimeoutStartSec=1800
