# systemd/mhde-monitor-model-perf.service
[Unit]
Description=MHDE monitor: rolling 7d precision vs walk-forward baseline
After=network.target

[Service]
Type=oneshot
User=jpcg
WorkingDirectory=/home/jpcg/MHDE
Environment=MHDE_DB_PATH=/home/jpcg/MHDE/data/mhde.duckdb
ExecStart=/home/jpcg/MHDE/venv/bin/python main.py monitor model-performance
StandardOutput=append:/home/jpcg/MHDE/data/logs/monitor_model_perf.log
StandardError=append:/home/jpcg/MHDE/data/logs/monitor_model_perf.log
