# systemd/mhde-monitor-streamlit-freshness.timer
[Unit]
Description=MHDE monitor: streamlit freshness — hourly

[Timer]
OnCalendar=*-*-* *:35:00
Persistent=true

[Install]
WantedBy=timers.target
