# systemd/mhde-crypto-export-predictions.timer
[Unit]
Description=MHDE crypto daily predictions export timer (06:15 UTC daily, 7 days/week)

[Timer]
OnCalendar=*-*-* 06:15:00
Persistent=true

[Install]
WantedBy=timers.target
