# systemd/mhde-monitor-phase0-calibration.timer
[Unit]
Description=MHDE monitor: Phase 0 calibration — weekly Sun 06:00 UTC

[Timer]
OnCalendar=Sun *-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
