# systemd/mhde-monitor-streamlit-freshness.service
[Unit]
Description=MHDE monitor: streamlit running-code freshness vs master
After=network.target

[Service]
Type=oneshot
User=jpcg
WorkingDirectory=/home/jpcg/MHDE
Environment=MHDE_DB_PATH=/home/jpcg/MHDE/data/mhde.duckdb
ExecStart=/home/jpcg/MHDE/venv/bin/python main.py monitor streamlit-freshness
StandardOutput=append:/home/jpcg/MHDE/data/logs/monitor_streamlit_freshness.log
StandardError=append:/home/jpcg/MHDE/data/logs/monitor_streamlit_freshness.log
