# systemd/mhde-monitor-model-perf.timer
[Unit]
Description=MHDE monitor: model performance daily

[Timer]
OnCalendar=*-*-* 13:15:00
Persistent=true

[Install]
WantedBy=timers.target
