# systemd/mhde-crypto-predict.service
[Unit]
Description=MHDE crypto daily prediction
After=network.target

[Service]
Type=oneshot
User=jpcg
WorkingDirectory=/home/jpcg/MHDE
Environment=MHDE_DB_PATH=/home/jpcg/MHDE/data/mhde.duckdb
ExecStart=/home/jpcg/MHDE/venv/bin/python main.py crypto backfill-prices
ExecStart=/home/jpcg/MHDE/venv/bin/python main.py crypto backfill-funding
ExecStart=/home/jpcg/MHDE/venv/bin/python main.py crypto backfill-oi
ExecStart=/home/jpcg/MHDE/venv/bin/python main.py crypto backfill-labels
ExecStart=/home/jpcg/MHDE/venv/bin/python main.py crypto backfill-features
ExecStart=/home/jpcg/MHDE/venv/bin/python main.py crypto predict
StandardOutput=append:/home/jpcg/MHDE/data/logs/crypto_predict.log
StandardError=append:/home/jpcg/MHDE/data/logs/crypto_predict.log
TimeoutStartSec=1800
