# systemd/mhde-crypto-retrain.timer
[Unit]
Description=MHDE crypto weekly retrain timer (Sun 20:00 UTC, before equity retrain at 21:30)

[Timer]
OnCalendar=Sun *-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
