# systemd/mhde-fx-predict.timer
[Unit]
Description=MHDE FX hourly prediction timer

[Timer]
OnCalendar=*-*-* *:05:00
Persistent=true

[Install]
WantedBy=timers.target
