# systemd/mhde-fx-retrain.timer
[Unit]
Description=MHDE FX weekly retrain timer (Sat 20:00 UTC, before any other writers)

[Timer]
OnCalendar=Sat *-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
