# systemd/mhde-fx-bot.service
[Unit]
Description=MHDE FX Telegram bot (long-polling)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jpcg
WorkingDirectory=/home/jpcg/MHDE
ExecStart=/home/jpcg/MHDE/venv/bin/python main.py fx bot
Environment=MHDE_DB_PATH=/home/jpcg/MHDE/data/mhde.duckdb
StandardOutput=append:/home/jpcg/MHDE/data/logs/fx_bot.log
StandardError=append:/home/jpcg/MHDE/data/logs/fx_bot.log
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=10

[Install]
WantedBy=multi-user.target
