# systemd/mhde-monitor-dashboard-synthetic.service
[Unit]
Description=MHDE monitor: dashboard synthetic E2E probe
After=network.target

[Service]
Type=oneshot
User=jpcg
WorkingDirectory=/home/jpcg/MHDE
Environment=MHDE_DB_PATH=/home/jpcg/MHDE/data/mhde.duckdb
ExecStart=/home/jpcg/MHDE/venv/bin/python main.py monitor dashboard-synthetic
StandardOutput=append:/home/jpcg/MHDE/data/logs/monitor_dashboard_synthetic.log
StandardError=append:/home/jpcg/MHDE/data/logs/monitor_dashboard_synthetic.log
