# systemd/mhde-capture-rest-collector.service
# /home/jpcg/.config/systemd/user/mhde-capture-rest-collector.service
[Unit]
Description=MHDE capture-core REST present-state collector (read-only, Binance public)
After=network-online.target
Wants=network-online.target

[Service]
# systemd --user unit (installed under ~/.config/systemd/user/): runs as the
# logging-in user, so NO User=/Group= line. Persistence relies on linger.
# Pre-flight fails fast if the venv is missing. The collector captures ONLY raw
# parquet under data/research/capture_core/ and never opens the production
# database. It self-paces under the live REST weight, coexisting with the WS
# capture service's depth SnapshotScheduler on the shared IP.
#
# BUILT-NOT-DEPLOYED: this file is tracked but not installed/enabled. R2 offload
# is DROPPED (ADR-035: forward-only). Resource caps (PR-3) express PRIORITY, not
# starvation: the engine wins contention; this collector is small and bounded.
Type=simple
WorkingDirectory=/home/jpcg/MHDE
ExecStartPre=/usr/bin/test -x /home/jpcg/MHDE/venv/bin/python
ExecStart=/home/jpcg/MHDE/venv/bin/python main.py crypto capture-rest-run
Restart=on-failure
RestartSec=5
# Give writers time to flush open parquet buffers on SIGTERM before SIGKILL.
TimeoutStopSec=30
# Resource caps (PR-3): a small REST poller, so a 1G ceiling is ample;
# CPUWeight/IOWeight < the system.slice default (100); OOM-first like the firehose.
MemoryMax=1G
CPUWeight=20
IOWeight=20
OOMScoreAdjust=800
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
