# systemd/mhde-crypto-export-predictions.timer
[Unit]
Description=MHDE crypto daily predictions export timer (00:40 UTC daily, 7 days/week — right after mhde-crypto-predict.service at 00:30)

[Timer]
OnCalendar=*-*-* 00:40:00
Persistent=true

[Install]
WantedBy=timers.target
