# .claude/local_scripts/daily_catalyst_queue.timer.example
[Unit]
Description=MHDE Daily Catalyst Queue Timer

[Timer]
OnCalendar=Mon..Fri *-*-* 23:15:00 UTC
# Note: does not skip US market holidays
Persistent=true

[Install]
WantedBy=timers.target
