# systemd/mhde-capture-core@.service.d/cpuset.conf
# /home/jpcg/.config/systemd/user/mhde-capture-core@.service.d/cpuset.conf
# ADR-039 gap 4 — pin ALL capture shards to the capture core band 4-7, disjoint from the
# engine band 0-3. A TEMPLATE drop-in: it applies to every instance @0..@7. N=8 shards share
# the 4-core band and the kernel balances the loops across it — shard-0 is heavier (it also
# owns the market-wide array streams), so a shared band beats rigid one-core-per-shard pinning.
#
# INERT until the cpuset-DELEGATION drop-in is installed: AllowedCPUs on a --user unit is a
# silent no-op while app.slice lacks the cpuset controller (ADR-039 §B verified host blocker).
# See systemd/system/user@.service.d/10-cpuset-delegate.conf (sudo). BUILT-NOT-DEPLOYED.
[Service]
AllowedCPUs=4-7
