# systemd/system/user@.service.d/10-cpuset-delegate.conf
# SYSTEM-level drop-in — DO NOT install with `systemctl --user`.
#
# Install target (operator, sudo):
#   sudo install -D -m 0644 \
#     systemd/system/user@.service.d/10-cpuset-delegate.conf \
#     /etc/systemd/system/user@.service.d/10-cpuset-delegate.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart user@$(id -u jpcg).service     # re-derives app.slice controllers
#                                                         # (or a full re-login)
#
# Why this is required (ADR-039 §B — the verified host blocker):
#   The capture shards run as systemd --user units, i.e. under app.slice. The user@.service
#   manager delegates only `cpu memory pids` to app.slice by default (NO cpuset), so the
#   AllowedCPUs=4-7 drop-ins on the capture units are a SILENT NO-OP — the directive shows in
#   the unit file but the kernel never programs cpuset.cpus. (system.slice DOES carry cpuset, so
#   the engine side enforces; that asymmetry is the whole problem.)
#
#   Extending user@.service's Delegate to include cpuset makes the user manager enable the
#   cpuset controller in app.slice's subtree, so the capture units' AllowedCPUs=4-7 ACTUALLY
#   pins them to cores 4-7, provably disjoint from the engine band 0-3.
#
# The ENGINE side (AllowedCPUs=0-3 on trading-engine-*) is the CROSS-REPO counterpart, owned by
# the crypto-trading-engine repo (its units live in /etc/systemd/system/). "Provably disjoint"
# holds only once ALL THREE are installed: this delegation + the capture band + the engine pin.
#
# Not auto-installed: system-level, needs sudo. The OPERATIONS.md runbook installs it before
# enabling mhde-capture.target.
[Service]
Delegate=cpu cpuset memory pids
