# systemd/mhde-capture-owner.service.d/cpuset.conf
# /home/jpcg/.config/systemd/user/mhde-capture-owner.service.d/cpuset.conf
# ADR-039 gap 4 — pin the snapshot-owner to the capture band 4-7, disjoint from the engine
# band 0-3. The owner is near-idle (request/response broker) so it shares the band with the
# shards rather than reserving a core.
#
# INERT until the cpuset-DELEGATION drop-in is installed (silent no-op on --user otherwise;
# ADR-039 §B). See systemd/system/user@.service.d/10-cpuset-delegate.conf. BUILT-NOT-DEPLOYED.
[Service]
AllowedCPUs=4-7
